# app.service
[Unit]
Description=Backend 102 app

[Service]
ExecStart=/usr/bin/java -jar /hello-world.jar

[Install]
WantedBy=multi-user.target
